# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/lib/systemd/system/alex.target
[Unit]
Description=The Alexandra target
After=alex-core.target
Requires=alex-core.target
AllowIsolate=yes

[Install]
WantedBy=alex-core.target
