# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/lib/systemd/system/alex-standalone.target
[Unit]
Description=Alexandra standalone base SW apps
After=alex-boot.target
Requires=alex-boot.target

[Install]
WantedBy=alex-boot.target
