# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/lib/systemd/system/alex-boot.target
[Unit]
Description=The alexandra boot target
After=network.target
Requires=network.target
AllowIsolate=yes

[Install]
WantedBy=network.target
RequiredBy=alex-core.target
