# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/lib/systemd/system/backendproxy.service
[Unit]
Description=The Alexandra Backend Proxy
After=alex-boot.target
Requires=alex-boot.target


[Service]
Type=notify
Restart=on-failure
ExecStart=/usr/bin/backendproxy




[Install]
WantedBy=alex-core.target


