# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/lib/systemd/system/service_tuner.service
[Unit]
Description=Service tuner
DefaultDependencies=no
After=ecuvard.service
Wants=ecuvard.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/service_tuner.sh

[Install]
WantedBy=alex-boot.target
