# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/usr/lib/cmake/cslclapi/in/systemd.service.in
[Unit]
Description=@S_UNIT_DESCRIPTION@
After=@S_UNIT_AFTER@
Requires=@S_UNIT_REQUIRES@
@S_UNIT_OTHER@

[Service]
Type=@S_SERVICE_TYPE@
Restart=@S_SERVICE_RESTART@
ExecStart=@S_SERVICE_EXEC@
@_S_EXEC_USER@
@S_SERVICE_OTHER@
@S_RESOURCECTL@

[Install]
WantedBy=@S_INSTALL_WANTEDBY@
@S_INSTALL_OTHER@

