# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-malc-cert.list
/lib/systemd/system/cert_app.service
/usr/bin/cert_app
/etc/alex/cert_app.conf
