# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-ecuvard.list
/lib/systemd/system/ecuvard.service
/usr/bin/ecuvard
