# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-modemd.list
/lib/systemd/system/modemd.service
/etc/systemd/network/usb0.network
/usr/bin/mdmctl.sh
/usr/bin/modemd
/etc/udev/rules.d/70-ublox-modem-at.rules
