# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-gnssmgr.list
/etc/gnssmgr/
/etc/gnssmgr/gnssmgr.conf
/usr/bin/gnssmgr
/lib/systemd/system/gnssmgr.service
