# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-malc-cert.postinst
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

if type systemctl >/dev/null 2>/dev/null; then
	systemctl $OPTS enable cert_app.service

	if [ -z "$D" -a "enable" = "enable" ]; then
		systemctl restart cert_app.service
	fi
fi
