# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-loctimed.list
/usr/lib/loctimed/tzdata/africa
/usr/lib/loctimed/tzdata/CONTRIBUTING
/etc/alex/loctimed.conf
/usr/lib/loctimed/tzdata/iso3166.tab
/usr/lib/loctimed/tzdata/systemv
/usr/lib/loctimed/tzdata/NEWS
/usr/lib/loctimed/tzdata/leapseconds.awk
/usr/lib/loctimed/
/usr/lib/loctimed/tzdata/calendars
/usr/lib/loctimed/tzdata/
/usr/lib/loctimed/tzdata/etcetera
/usr/lib/loctimed/tzdata/australasia
/usr/lib/loctimed/tzdata/europe
/usr/lib/loctimed/tzdata/southamerica
/usr/lib/loctimed/tzdata/leapseconds
/usr/lib/loctimed/tzdata/leap-seconds.list
/usr/lib/loctimed/tzdata/zone1970.tab
/usr/lib/loctimed/tzdata/Makefile
/usr/lib/loctimed/tzdata/LICENSE
/usr/lib/loctimed/tzdata/zone.tab
/usr/lib/loctimed/tzdata/antarctica
/usr/lib/loctimed/tzdata/asia
/usr/lib/loctimed/zddatabase.bin
/usr/bin/loctimed
/usr/lib/loctimed/tzdata/checktab.awk
/lib/systemd/system/loctimed.service
/usr/lib/loctimed/tzdata/zishrink.awk
/usr/lib/loctimed/tzdata/yearistype.sh
/usr/lib/loctimed/tzdata/version
/usr/lib/loctimed/tzdata/factory
/usr/lib/loctimed/tzdata/README
/usr/lib/loctimed/tzdata/northamerica
/usr/lib/loctimed/tzdata/checklinks.awk
/usr/lib/loctimed/tzdata/theory.html
/usr/lib/loctimed/tzdata/backward
/usr/lib/loctimed/tzdata/zoneinfo2tdf.pl
/usr/lib/loctimed/tzdata/backzone
/usr/lib/loctimed/tzdata/pacificnew
/usr/lib/loctimed/tzdata/ziguard.awk
