# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-vpnctld.list
/usr/bin/vpnctld
/lib/systemd/system/vpnctld.service
