# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-udsifd.list
/usr/bin/udsifd
/lib/systemd/system/udsifd.service
/etc/alex/udsifd.conf
