# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-malc-systime.list
/etc/systimemgr/
/etc/systimemgr/systimemgr.conf
/lib/systemd/system/systimemgr.service
/usr/bin/systimemgr
