# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-pgnreqd.list
/etc/alex/pgnreqd.conf
/usr/bin/pgnreqd
/lib/systemd/system/pgnreqd.service
