# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-devtools.list
/lib/systemd/system/ws_server.service
/usr/bin/logctl
/usr/bin/udsifc
/usr/bin/ws_server
/usr/bin/beproxyc
/usr/bin/evild
/usr/bin/aclifc
/usr/bin/smoketest.sh
/lib/systemd/system/ws_server@.service
/usr/bin/emmcc
/lib/systemd/system/container-test@.service
/usr/bin/cert-test
/usr/bin/bootctl
/etc/alex/container-test.conf
/usr/bin/container-test
/usr/bin/aleavesdrop
/usr/bin/sysparamsc
/lib/systemd/system/container-test.service
/etc/alex/ws_server.conf
/usr/bin/rtsigc
/usr/bin/message_sender
/usr/bin/acctl
/usr/bin/maintrc
