# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-cslc-servicetuner.list
/usr/sbin/service_tuner.sh
/lib/systemd/system/service_tuner.service
