# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/alex-malc-beproxy.list
/lib/systemd/system/backendproxy.service
/var/spool/backendproxy/mq/
/etc/alex/backendproxy.conf
/usr/bin/backendproxy
/var/spool/backendproxy/
