# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/bluez5.list
/etc/bluetooth/network.conf
/etc/bluetooth/input.conf
/usr/bin/mpris-proxy
/usr/share/dbus-1/system-services/org.bluez.service
/usr/lib/libbluetooth.so.3.18.14
/usr/bin/hciconfig
/usr/bin/hcidump
/lib/udev/hid2hci
/lib/systemd/system/bluetooth.service
/usr/bin/ciptool
/usr/bin/l2ping
/usr/bin/hcitool
/usr/lib/bluez5/bluetooth/
/usr/bin/bluemoon
/usr/bin/sdptool
/usr/bin/hex2hcd
/usr/bin/bluetoothctl
/usr/bin/hciattach
/usr/bin/bccmd
/usr/bin/l2test
/usr/lib/bluez5/
/lib/udev/rules.d/97-hid2hci.rules
/usr/lib/libbluetooth.so.3
/etc/dbus-1/system.d/bluetooth.conf
/usr/bin/rctest
/etc/bluetooth/
/usr/bin/btmon
/usr/lib/bluez5/bluetooth/bluetoothd
/usr/bin/rfcomm
