# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.d/
/etc/dnsmasq.conf
/usr/bin/dnsmasq
/usr/bin/dhcp_release
/lib/systemd/system/dnsmasq.service
