# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/emmy-w1-driver-tools.list
/usr/sbin/reset_wifi.sh
/lib/systemd/system/configure_wifi.service
/lib/systemd/system/reset_wifi.service
/usr/sbin/configure_wifi.sh
