# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/hostapd.list
/usr/sbin/hostapd
/lib/systemd/system/hostapd.service
/usr/sbin/hostapd_cli
/etc/hostapd.conf
