# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/mosquitto.list
/etc/mosquitto/mosquitto.conf
/usr/bin/mosquitto_passwd
/usr/sbin/mosquitto
/etc/mosquitto/
/lib/systemd/system/mosquitto.service
