# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/openssh-sshd.list
/lib/
/usr/
/lib/systemd/system/sshd.socket
/lib/systemd/
/etc/
/lib/systemd/system/
/lib/systemd/system/sshdgenkeys.service
/usr/sbin/sshd
/etc/ssh/sshd_config_readonly
/etc/ssh/moduli
/
/lib/systemd/system/sshd@.service
/etc/ssh/
/usr/sbin/
/etc/ssh/sshd_config
/etc/default/volatiles/99_sshd
/etc/default/volatiles/
/etc/default/
