# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/sre-systemd-target.postinst
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

if type systemctl >/dev/null 2>/dev/null; then
	systemctl $OPTS enable  alex-boot.target alex-core.target alex-standalone.target alex.target

	if [ -z "$D" -a "enable" = "enable" ]; then
		systemctl restart  alex-boot.target alex-core.target alex-standalone.target alex.target
	fi
fi
