# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/sre-ipcd.list
/usr/bin/ipcd
/etc/alex/ipcd.conf
/lib/systemd/system/ipcd.service
