# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/sre-spider.list
/usr/bin/spider
/lib/systemd/system-shutdown/spider
/var/spool/spider/
/lib/systemd/system/spider.service
/etc/alex/spider.conf
/var/spool/spider/updates/
