# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/sre-check-ethernet-link-status.list
/lib/systemd/system/check_ethernet_link_status.service
/usr/sbin/check_ethernet_link_status.sh
