# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/sre-systemd-target.list
/lib/systemd/system/alex-boot.target
/lib/systemd/system/alex-standalone.target
/lib/systemd/system/alex.target
/lib/systemd/system/alex-core.target
