# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/sre-usblog.list
/usr/share/log.public.pem
/usr/bin/extract_log.sh
/lib/systemd/system/log_extractor.service
/etc/alex/partition.args
