# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/sre-demo-app.list
/lib/systemd/system/sre-demo-app.service
/usr/bin/sre-demo-app
