# SDK/sysroots/cortexa9hf-vfp-neon-poky-linux-gnueabi/var/lib/opkg/info/sre-usblog-move-tpa-logs.list
/lib/systemd/system/move_tpa_logs.service
/usr/bin/move_tpa_logs.sh
